# server/com.abe.open-brain-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.abe.open-brain-api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/mikeingersoll/abe-command-center/server/index.mjs</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/open-brain-api.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/open-brain-api.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>OPEN_BRAIN_TOKEN</key>
        <string>abe-open-brain-2026</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
